# deploy/nas-fm.service
[Unit]
Description=Self-hosted NAS File Manager (React + Hono)
After=network.target

[Service]
User=<あなたのユーザー名>
Group=nas
UMask=0002
WorkingDirectory=/opt/nas-fm
Environment=NAS_ROOT=/srv/nas/share
Environment=PORT=8080
# AUTH_SECRET / AUTH_PASSWORD_HASH は秘密値のため Git 管理外の EnvironmentFile に分離する。
# /opt/nas-fm/nas-fm.env は chmod 600 にし、root と Group= のユーザーのみ読めるようにすること。
# AUTH_SECRET はランダムな長い文字列。
# AUTH_PASSWORD_HASH は `npx tsx apps/server/scripts/hash-password.ts <password>` の出力を使う。
EnvironmentFile=/opt/nas-fm/nas-fm.env
ExecStart=/usr/bin/node /opt/nas-fm/server.js
Restart=on-failure

[Install]
WantedBy=multi-user.target
